# profile.d/codeplag
# codeplag autocomplete initialization script (bash)
eval " $(register-python-argcomplete codeplag)"
